# daily-article.service.template
[Unit]
Description=Serve MDN articles.

[Service]
Type=simple
ExecStart=#start_script#

[Install]
WantedBy=multi-user.target
